# mihomo/rc.d/mihomo
#!/bin/sh
#
# PROVIDE: mihomo
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="mihomo"
rcvar="${name}_enable"
desc="mihomo service"

pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"
procname="/usr/local/bin/mihomo"
command="/usr/sbin/daemon"

# 必须存在的文件
required_files="/usr/local/bin/mihomo /usr/local/etc/mihomo/config.yaml"

# 自定义命令钩子
start_precmd="mihomo_prestart"
start_postcmd="mihomo_poststart"
stop_cmd="mihomo_stop"
status_cmd="mihomo_status"
extra_commands="reload"
reload_cmd="mihomo_reload"

# 环境变量 & TUN 设备名
mihomo_env="IS_DAEMON=1"
tun_device="tun_3000"

# daemon 启动参数（核心修复：环境变量 + 正确命令）
command_args="-P ${pidfile} -f -o ${logfile} -t ${name} env ${mihomo_env} ${procname} -d /usr/local/etc/mihomo/"

# 默认不启用
: ${mihomo_enable:="NO"}

# 启动前检查
mihomo_prestart() {
    touch "${logfile}"

    # 清理无效 PID 文件
    if [ -f "${pidfile}" ]; then
        if kill -0 "$(cat "${pidfile}")" 2>/dev/null; then
            echo "${name} is already running."
            return 1
        fi
        rm -f "${pidfile}"
    fi

    # 禁止多开
    if pgrep -fx "${procname} -d /usr/local/etc/mihomo/" >/dev/null 2>&1; then
        echo "A mihomo process is already running, abort."
        return 1
    fi

    return 0
}

# 启动后验证
mihomo_poststart() {
    sleep 1
    if [ -s "${pidfile}" ] && kill -0 "$(cat "${pidfile}")" 2>/dev/null; then
        mihomo_pid=$(pgrep -P $(cat ${pidfile}) ${procname} 2>/dev/null | head -1)
        if [ -n "${mihomo_pid}" ]; then
            echo "✅ ${name} started (daemon: $(cat ${pidfile}), mihomo: ${mihomo_pid})"
            return 0
        fi
    fi

    echo "❌ ${name} failed to start, check log: ${logfile}"
    return 1
}

# 停止服务（含 TUN 清理）
mihomo_stop() {
    local daemon_pid mihomo_pids extra_pids i

    echo "🛑 Stopping ${name}..."

    # 1. 停止 daemon + mihomo
    if [ -f "${pidfile}" ]; then
        daemon_pid=$(cat "${pidfile}")
        mihomo_pids=$(pgrep -P ${daemon_pid} ${procname} 2>/dev/null)

        [ -n "${mihomo_pids}" ] && kill -TERM ${mihomo_pids} 2>/dev/null
        kill -TERM ${daemon_pid} 2>/dev/null

        # 等待退出
        for i in 1 2 3 4 5; do
            kill -0 ${daemon_pid} 2>/dev/null || break
            sleep 1
        done

        # 强杀
        kill -KILL ${mihomo_pids} ${daemon_pid} 2>/dev/null
    fi

    # 2. 清理孤儿进程
    extra_pids=$(pgrep -fx "${procname} -d /usr/local/etc/mihomo/" 2>/dev/null)
    [ -n "${extra_pids}" ] && kill -KILL ${extra_pids} 2>/dev/null

    # 3. 清理 TUN 设备
    if [ -n "${tun_device}" ] && ifconfig | grep -q "${tun_device}"; then
        ifconfig ${tun_device} destroy 2>/dev/null
        echo "✅ Cleaned up TUN device: ${tun_device}"
    fi

    # 4. 清理 PID 文件
    rm -f "${pidfile}"
    echo "✅ ${name} stopped."
    return 0
}

# 状态查询
mihomo_status() {
    if [ -f "${pidfile}" ] && kill -0 "$(cat ${pidfile})" 2>/dev/null; then
        local dpid=$(cat ${pidfile})
        local mpid=$(pgrep -P ${dpid} ${procname} 2>/dev/null | head -1)
        if [ -n "${mpid}" ]; then
            echo "✅ ${name} is running (daemon: ${dpid}, mihomo: ${mpid})"
            return 0
        fi
    fi

    local orphan=$(pgrep -fx "${procname} -d /usr/local/etc/mihomo/" 2>/dev/null | head -1)
    if [ -n "${orphan}" ]; then
        echo "⚠️ ${name} running (orphan pid: ${orphan}), please run: service ${name} stop"
        return 1
    fi

    echo "❌ ${name} is not running."
    return 1
}

# 热重载配置
mihomo_reload() {
    if [ -f "${pidfile}" ]; then
        local dpid=$(cat ${pidfile})
        local mpid=$(pgrep -P ${dpid} ${procname} 2>/dev/null | head -1)
        if [ -n "${mpid}" ]; then
            kill -HUP ${mpid}
            echo "🔄 ${name} config reloaded."
            return 0
        fi
    fi

    echo "❌ ${name} is not running."
    return 1
}

# 加载 rc.conf 配置
load_rc_config "${name}"

# 命令行手动执行时临时允许启动
if [ -n "$1" ]; then
    eval "${rcvar}=YES"
fi

run_rc_command "$1"