# rc.d/clash
#!/bin/sh
#
# PROVIDE: clash
# REQUIRE: NETWORKING SERVERS
# KEYWORD: shutdown

. /etc/rc.subr

name="clash"
rcvar=clash_enable
pidfile="/var/run/${name}.pid"
log_file="/var/log/${name}.log"

# 你的Clash路径和配置目录（和WebUI完全对齐）
clash_bin="/usr/local/bin/clash"
clash_conf_dir="/usr/local/etc/clash"

# 开机默认不自启
clash_enable=${clash_enable:-"NO"}

# 自定义启动/停止/状态命令
start_cmd="my_clash_start"
stop_cmd="my_clash_stop"
status_cmd="my_clash_status"

# 启动函数
my_clash_start() {
    echo "正在启动 Clash..."

    # 清理残留进程和PID文件
    rm -f "${pidfile}"
    pkill -9 "${name}" >/dev/null 2>&1
    sleep 1

    # 后台启动Clash（和你手动执行的命令完全一致）
    cd "${clash_conf_dir}" || exit 1
    nohup "${clash_bin}" -d "${clash_conf_dir}" > "${log_file}" 2>&1 &
    
    # 记录进程PID
    echo $! > "${pidfile}"
    
    # 验证启动结果
    sleep 2
    if [ -f "${pidfile}" ] && kill -0 $(cat "${pidfile}") 2>/dev/null; then
        echo "✅ Clash 启动成功！PID: $(cat "${pidfile}")"
        return 0
    else
        echo "❌ Clash 启动失败，请查看日志: ${log_file}"
        rm -f "${pidfile}"
        return 1
    fi
}

# 停止函数
my_clash_stop() {
    echo "正在停止 Clash..."
    if [ -f "${pidfile}" ]; then
        kill -9 $(cat "${pidfile}") 2>/dev/null
        rm -f "${pidfile}"
    fi
    pkill -9 "${name}" >/dev/null 2>&1
    echo "✅ Clash 已停止。"
}

# 状态检查函数
my_clash_status() {
    if [ -f "${pidfile}" ] && kill -0 $(cat "${pidfile}") 2>/dev/null; then
        echo "✅ Clash 正在运行 (PID: $(cat "${pidfile}"))"
        return 0
    else
        echo "❌ Clash 未运行"
        return 1
    fi
}

# 加载系统配置
load_rc_config $name

# 【关键修复】强制允许手动执行start命令，不影响开机不自启
# 开机时系统只有clash_enable=YES才会调用本脚本，因此这里强制设为YES仅影响手动执行
clash_enable="YES"

# 执行命令
run_rc_command "$1"